# rtl-fm-radio.service
[Unit]
Description=RTL-SDR Web Radio Server
Documentation=https://github.com/your-repo/rtl_fm_python_webgui
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/rtl_fm_python_webgui
ExecStart=/home/pi/rtl_fm_python_webgui/rtl_fm_python_web.py -M wbfm -f 101.1M -
Restart=on-failure
RestartSec=5

# Give the RTL-SDR time to initialize
ExecStartPre=/bin/sleep 2

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rtl-fm-radio

[Install]
WantedBy=multi-user.target
